# majordomus-control/rootfs/etc/services.d/majordomus/finish
#!/usr/bin/execlineb -S0

# Bring down the container when the service exits
s6-svscanctl -t /var/run/s6/services
