# majordomus-control/rootfs/etc/services.d/majordomus/finish
#!/usr/bin/execlineb -P

# Service exited — s6 will handle restart automatically for longrun services
